# provisioning/wifi-connect/qdrn-wifi-connect.service
# Brings up the QDRN-Radar-Setup captive portal when the Pi has no internet.
# wifi-connect exits once the friend selects a network + it connects, after
# which normal WiFi resumes. It re-arms on next boot if connectivity is lost.
[Unit]
Description=QDRN Radar WiFi captive portal (wifi-connect)
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
# Only start the portal if we are NOT already online.
ExecStartPre=/bin/sh -c 'nm-online -q -t 20 && exit 1 || exit 0'
ExecStart=/usr/local/sbin/wifi-connect --portal-ssid "QDRN-Radar-Setup"
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
