# provisioning/wifi-connect/qdrn-wifi-connect.service
# Brings up the QDRN-Radar-Setup captive portal whenever the Pi is offline.
# The wrapper exits 0 when online; systemd restarts it on a delay so the Pi
# effectively polls and re-opens the portal if WiFi is removed/forgotten later.
# When offline, wifi-connect runs the hotspot and blocks until the user picks a
# network, then exits — and the poll cycle resumes.
[Unit]
Description=QDRN Radar WiFi captive portal (wifi-connect)
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/qdrn-portal.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
