# provisioning/wifi-connect/qdrn-wifi-connect.service
# Brings up the QDRN-Radar-Setup captive portal at boot, but only if offline.
# The wrapper exits 0 when already connected, so there's no restart loop on a
# Pi that's working normally. wifi-connect exits once the user picks a network.
[Unit]
Description=QDRN Radar WiFi captive portal (wifi-connect)
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/qdrn-portal.sh
Restart=no

[Install]
WantedBy=multi-user.target
